# overlay/elvees/usr/lib/systemd/system/canfdechoserver@.service
[Unit]
Description=CAN full-duplex echo server

[Service]
ExecStart=/bin/canfdtest %i

[Install]
WantedBy=multi-user.target
